# magnetodb/debian/magnetodb.prerm
#!/bin/sh

set -e

case $1 in
	remove|purge)
		if [ -x /etc/init.d/magnetodb ]; then
			invoke-rc.d magnetodb stop || true
		fi
		;;
esac

#DEBHELPER#
